# crontab
30 4 * * * root /opt/backup/mongodb-backup.sh > /dev/null 2>&1
0 5 * * * root /opt/backup/elastic_snapshot.sh > /dev/null 2>&1
0 6 * * * root /opt/backup/remove_backup.sh > /dev/null > 2>&1
